# yongbeom-hetzner-auction-ubuntu-2404/vaultwarden/backup/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

20 4 * * * /usr/local/bin/backup.sh
