# yongbeom-hetzner-auction-ubuntu-2404/keycloak/backup/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

25 3 * * * /usr/local/bin/backup.sh
